# s6/services/.s6-svscan/finish
#!/bin/sh

echo "Clean exit from Envoymon"
